# overrides/supervisord/development/django/celerybeat.conf
; ========================
;  celery beat supervisor
; ========================

; the name of your supervisord program
[program:hubcelerybeat]

; Set full path to celery program if using virtualenv
command=/usr/local/bin/celery beat --app=hub.celery --loglevel=INFO --pidfile=/home/celery/celerybeat.pid --schedule=/home/celery/celerybeat-schedule

; The directory to our Django project
directory=/src/

; If supervisord is run as the root user, switch users to this UNIX user account
; before doing any processing.
user=celery

; Supervisor will start as many instances of this program as named by numprocs
numprocs=1

; we are in a development environment and want to see our stdout directly
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

; May be one of false, unexpected, or true. If false, the process will never
; be autorestarted. If unexpected, the process will be restart when the program
; exits with an exit code that is not one of the exit codes associated with this
; process’ configuration (see exitcodes). If true, the process will be
; unconditionally restarted when it exits, without regard to its exit code.
autorestart=true

; The total number of seconds which the program needs to stay running after
; a startup to consider the start successful.
startsecs=10
